# systemd/.config/systemd/user/rdm.service
[Unit]
Description=RTags deamon

[Service]
ExecStart=%h/Code/Cpp/rtags/bin/rdm

[Install]
WantedBy=default.target
